# systemd/collabhost.user.service
# Collabhost -- systemd user-scope unit template
#
# Copy this file to ~/.config/systemd/user/collabhost.service, then enable:
#
#   loginctl enable-linger "$USER"          # keep the user manager alive after logout
#   systemctl --user daemon-reload
#   systemctl --user enable --now collabhost.service
#
# Logs:
#
#   journalctl --user -u collabhost -f
#
# This is the user-scope template (phase 2 of card #230). A system-scope variant
# (collabhost.system.service) is planned for phase 3 -- it will run under a
# dedicated `collabhost` system user with AmbientCapabilities=CAP_NET_BIND_SERVICE,
# bound to /opt/collabhost and /var/lib/collabhost. Ship date TBD.
#
# Edit the marked lines below to match your install path. The defaults assume the
# install.sh layout ($HOME/.collabhost/bin).

[Unit]
Description=Collabhost
Documentation=https://github.com/MrBildo/collabhost
After=network-online.target
Wants=network-online.target

[Service]
Type=simple

# --- EDIT THESE TWO LINES if you installed to a non-default path ---
WorkingDirectory=%h/.collabhost/bin
ExecStart=%h/.collabhost/bin/collabhost
# -------------------------------------------------------------------

# Restart on crash. RestartSec gives the host a moment to flush stderr / crash
# logs before systemd respawns. StartLimitBurst=5 / StartLimitIntervalSec=60
# stop a tight crash loop from spinning forever.
Restart=on-failure
RestartSec=5

# Operator-relevant environment overrides. Uncomment and set as needed -- these
# layer on top of appsettings.json (env > appsettings > built-in default). See
# INSTALL.md section 5 for the full reference.
#
# Environment="COLLABHOST_DATA_PATH=%h/collabhost-data"
# Environment="COLLABHOST_PROXY_BASE_DOMAIN=collabhost.lan"
# Environment="COLLABHOST_ADMIN_KEY=01JABCDEFGHJKMNPQRSTVWXYZ"
#
# DNS-01 ACME issuer (real Let's Encrypt certs on a public domain). Required when
# Proxy:DnsProvider is set in appsettings.json. Replace with your provider's token.
#
# Environment="CLOUDFLARE_API_TOKEN=replace-me"

# Conservative process hygiene. Tighten further if your deployment shape allows
# (PrivateTmp, ProtectSystem, NoNewPrivileges, etc. are standard candidates -- not
# enabled here because they interact with the bundled Caddy's privileged-port
# binding and the data/ write path; tighten when you've validated against your
# own layout).

[Install]
WantedBy=default.target
